# audio-pi.service
[Unit]
Description=Audio Pi Websystem
After=network.target

[Service]
Type=simple
User=pi
Group=audio
WorkingDirectory=/opt/Audio-Pi-Websystem
ExecStartPre=/bin/sleep 10
Environment=FLASK_DEBUG=0
Environment=FLASK_SECRET_KEY=
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/opt/Audio-Pi-Websystem/venv/bin/python /opt/Audio-Pi-Websystem/app.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
