# audio-pi.service
[Unit]
Description=Audio Pi Websystem
After=network.target

[Service]
Type=simple
User=pi
Group=pi
RuntimeDirectory=audio-pi
RuntimeDirectoryMode=0700
WorkingDirectory=/opt/Audio-Pi-Websystem
ExecStartPre=/bin/sleep 10
Environment=FLASK_DEBUG=0
Environment="FLASK_SECRET_KEY=__CHANGE_ME__"
Environment=FLASK_PORT=80
Environment=XDG_RUNTIME_DIR=/run/audio-pi
Environment=AUDIO_PI_DISABLE_SUDO=0
ExecStart=/opt/Audio-Pi-Websystem/venv/bin/gunicorn --config /opt/Audio-Pi-Websystem/gunicorn.conf.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillSignal=SIGQUIT
TimeoutStartSec=60
TimeoutStopSec=45
Restart=on-failure
RestartSec=5
NoNewPrivileges=no
RestrictSUIDSGID=no
ProtectSystem=strict
ReadWritePaths=/opt/Audio-Pi-Websystem
PrivateTmp=yes
ProtectControlGroups=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
LockPersonality=yes
RestrictNamespaces=yes
RestrictRealtime=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK AF_BLUETOOTH
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_BOOT CAP_SYS_TIME CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_BIND_SERVICE CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_BOOT CAP_SYS_TIME CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
